# config/mounts/media-Synology-Linux.automount
[Unit]
Description=Automount Synology NAS through SMB protocol - Linux shared folder

[Automount]
Where=/media/Synology/Linux
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
