# config/mounts/media-Synology-Backup.automount
[Unit]
Description=Automount Synology NAS through SMB protocol - Backup shared folder

[Automount]
Where=/media/Synology/Backup
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
