# config/mounts/media-Synology-Fotos.automount
[Unit]
Description=Automount Synology NAS through SMB protocol - Fotos shared folder

[Automount]
Where=/media/Synology/Fotos
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
