# config/mounts/media-Synology-Docker.automount
[Unit]
Description=Automount Synology NAS through SMB protocol - Docker shared folder

[Automount]
Where=/media/Synology/Docker
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
