# config/mounts/media-Synology-Usb.automount
[Unit]
Description=Automount Synology NAS through SMB protocol - USB shared folder

[Automount]
Where=/media/Synology/Usb
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
