# config/mounts/media-Synology-Usb3.automount
[Unit]
Description=Automount Synology NAS through SMB protocol - USB3 Shared Folder
Where=/media/Synology/Usb3
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
