# config/mounts/media-Synology-Documentos.automount
[Unit]
Description=Automount Synology NAS through SMB protocol - Documentos shared folder

[Automount]
Where=/media/Synology/Documentos
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
