# config/mounts/media-Synology-Misc.automount
[Unit]
Description=Automount Synology NAS through SMB protocol - Misc shared folder

[Automount]
Where=/media/Synology/Misc
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
