# config/mounts/media-Synology-Biblioteca.automount
[Unit]
Description=Automount Synology NAS through SMB protocol - Biblioteca shared folder

[Automount]
Where=/media/Synology/Biblioteca
TimeoutIdleSec=10

[Install]
WantedBy=multi-user.target
